# .config/autostart/sshd.desktop
[Desktop Entry]
Type=Application
Name=sshd-desktop
Exec=~/.config/ssh/sshd.sh
